# mikromon.service
[Unit]
Description=MikroMon Dashboard Backend
After=network.target

[Service]
User=root
WorkingDirectory=/home/sairo/Antigravity
ExecStart=/home/sairo/Antigravity/venv/bin/python /home/sairo/Antigravity/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
